# templates/default/kube-proxy.service.erb
[Unit]
Description=Kubernetes Kube-Proxy Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=network.target

[Service]
ExecStart=<%= node['kubernetes']['root'] %>/kubernetes/server/bin/kube-proxy <%= node['kubernetes']['args']['proxy'] %>
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
